# home/dot_config/systemd/user/dconf-drift.service
[Unit]
Description=Check for dconf drift
Documentation=man:dconf(1)

[Service]
ExecStart=%h/.local/bin/dconf-drift
