# home/dot_config/systemd/user/flatpak-drift.service
[Unit]
Description=Check for Flatpak drift
Documentation=man:flatpak(1)

[Service]
ExecStart=%h/.local/bin/flatpak-drift
